# roles/redmine-pgbouncer-haproxy/files/haproxy/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid

[unix_http_server]
file=/var/run/supervisor/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:consul]
command=/usr/local/bin/consul agent -config-dir=/consul/config -data-dir=/consul/data
stdout_logfile=/var/log/supervisor/consul.log
stderr_logfile=/var/log/supervisor/consul.err.log
priority=1
startsecs=3
startretries=3
autorestart=true

[program:haproxy]
command=/usr/local/sbin/haproxy -W -db -f /usr/local/etc/haproxy/haproxy.cfg -p /var/run/haproxy/haproxy.pid
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
priority=2
startsecs=3
startretries=3
autorestart=true

[program:delayed-startup]
command=/delayed-startup.sh
stdout_logfile=/var/log/supervisor/delayed-startup.log
stderr_logfile=/var/log/supervisor/delayed-startup.err.log
priority=3
startsecs=0
startretries=1
autorestart=false

[program:consul-template]
command=/usr/local/bin/consul-template -config /consul-template/config/consul-template.hcl
stdout_logfile=/var/log/supervisor/consul-template.log
stderr_logfile=/var/log/supervisor/consul-template.err.log
autostart=false
autorestart=true
