# ansible/roles/redmine-pgbouncer-haproxy/files/pgbouncer/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid

[program:consul]
command=/usr/local/bin/consul agent -data-dir=/consul/data -config-dir=/consul/config
stdout_logfile=/var/log/supervisor/consul.log
stderr_logfile=/var/log/supervisor/consul.err.log
priority=1
startsecs=3
startretries=3
autorestart=true

[program:pgbouncer]
command=/opt/bitnami/scripts/pgbouncer/entrypoint.sh /opt/bitnami/scripts/pgbouncer/run.sh
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
priority=2
startsecs=3
startretries=3
autorestart=true
